# scripts/packages/hbm.service
[Unit]
Description=HBM Docker Authorization Plugin
Documentation=https://github.com/jonasbroms/hbm
Before=docker.service
After=network.target hbm.socket
Requires=hbm.socket
Wants=docker.service

[Service]
Type=simple
ExecStartPre=-/usr/local/sbin/hbm init
ExecStart=/usr/local/sbin/hbm server
Restart=on-failure
RestartSec=10s

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ReadWritePaths=/var/lib/hbm /run/docker/plugins /etc/docker/plugins

[Install]
WantedBy=multi-user.target
