# scripts/packages/hbm.socket
[Unit]
Description=HBM Docker Authorization Plugin Socket

[Socket]
ListenStream=/run/docker/plugins/hbm.sock
SocketMode=0660

[Install]
WantedBy=sockets.target
